# systemd/system/arpsettings.service
[Unit]
Description=ARP Settings Service

[Service]
ExecStart=/etc/scripts/arp_settings.sh

[Install]
WantedBy=multi-user.target
